# supervisor-app.conf
[program:bothub-port-80]
command = /home/app/webapp/supervisor_command.sh