# supervisor-app.conf
[program:bothub-port-80]
command = /home/app/webapp/bothub/app/server.py 80