# supervisor-app.conf
[program:bothub-port-8000]
command = /home/app/webapp/bothub/app/server.py 8000

[program:nginx-app]
command = /usr/sbin/nginx